# systemd/.config/systemd/user/restic-backup.timer
[Unit]
Description=Backup with restic daily
[Timer]
OnCalendar=*-*-* 19:00:00
Persistent=True
[Install]
WantedBy=timers.target